# crates/bin/cf/docs/ddns.service
[Service]
ExecStart = %h/.local/bin/cf dns update
TimeoutSec = 1min

[Unit]
Description = DDNS
