# crates/bin/cf/docs/ddns.timer
[Install]
WantedBy = timers.target

[Timer]
OnStartupSec = 0
OnUnitActiveSec = 1min

[Unit]
Description = DDNS
